# services/service2/run
#!/usr/bin/with-contenv sh
echo >&2 "starting service2"
exec sleep 99999999
